# 1.1-udev_rule/80-test.rules
SUBSYSTEM=="block", ACTION=="add", ENV{DEVTYPE}=="disk", ATTRS{serial}=="000AEB91BC94A920557001AE", RUN+="/usr/bin/device_added.sh"

